# .config/systemd/user/cros-garcon.service.d/cros-garcon-override.conf
[Service]
ExecStart=
ExecStart=/opt/google/cros-containers/bin/garcon --allow_any_user --server
